# packages/exoctl_amd64/usr/lib/systemd/system/infinity-engine.service
[Unit]
Description=Infinity Engine
After=network.target

[Service]
Type=simple
ExecStart=/opt/infinity-engine/App
Environment="EXOCTLDIR=/opt/infinity-engine/"

[Install]
WantedBy=multi-user.target